# deploy/lake-prune.timer
# /etc/systemd/system/lake-prune.timer
[Unit]
Description=Daily lake prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
